# deploy/supervisor/celery.conf
; Production Supervisor config for Celery (FDMS multi-tenant Ping + other tasks).
; Install: sudo cp deploy/supervisor/celery.conf /etc/supervisor/conf.d/celery.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start celery_worker celery_beat
;
; Replace /var/www/fdms with your app root. Ensure redis is running.

[program:celery_worker]
command=/var/www/fdms/venv/bin/celery -A fdms_project worker -l info --concurrency=4
directory=/var/www/fdms
user=www-data
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
redirect_stderr=true
stdout_logfile=/var/log/celery/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=DJANGO_SETTINGS_MODULE="fdms_project.settings",CELERY_BROKER_URL="redis://127.0.0.1:6379/0"

[program:celery_beat]
command=/var/www/fdms/venv/bin/celery -A fdms_project beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/var/www/fdms
user=www-data
numprocs=1
autostart=true
autorestart=true
startsecs=10
redirect_stderr=true
stdout_logfile=/var/log/celery/beat.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=DJANGO_SETTINGS_MODULE="fdms_project.settings",CELERY_BROKER_URL="redis://127.0.0.1:6379/0"
