# deploy/supervisor/fdms.conf
; Supervisor config for FDMS SaaS (Ubuntu 22.04 / DigitalOcean).
; Install: sudo cp deploy/supervisor/fdms.conf /etc/supervisor/conf.d/fdms.conf
; Then: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start all
;
; Replace /var/www/fdms with your app root.
; Replace www-data with your run user.
; Set environment in [program:...] or via /etc/supervisor/conf.d/fdms.env (not committed).

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; --- FDMS application programs ---

[program:fdms_gunicorn]
command=/var/www/fdms/venv/bin/gunicorn fdms_project.wsgi:application --config /var/www/fdms/gunicorn.conf.py
directory=/var/www/fdms
user=www-data
autostart=true
autorestart=true
startretries=3
stopwaitsecs=120
environment=DJANGO_SETTINGS_MODULE="fdms_project.settings_production",PATH="/var/www/fdms/venv/bin:%(ENV_PATH)s"
stdout_logfile=/var/log/fdms/gunicorn_access.log
stderr_logfile=/var/log/fdms/gunicorn_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
redirect_stderr=false

[program:fdms_celery_worker]
command=/var/www/fdms/venv/bin/celery -A fdms_project worker -l info --concurrency=4
directory=/var/www/fdms
user=www-data
autostart=true
autorestart=true
startretries=3
stopwaitsecs=60
environment=DJANGO_SETTINGS_MODULE="fdms_project.settings_production",CELERY_BROKER_URL="redis://127.0.0.1:6379/0",CELERY_RESULT_BACKEND="redis://127.0.0.1:6379/0",PATH="/var/www/fdms/venv/bin:%(ENV_PATH)s"
stdout_logfile=/var/log/fdms/celery_worker.log
stderr_logfile=/var/log/fdms/celery_worker_error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
redirect_stderr=false

[program:fdms_celery_beat]
command=/var/www/fdms/venv/bin/celery -A fdms_project beat -l info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/var/www/fdms
user=www-data
autostart=true
autorestart=true
startretries=3
stopwaitsecs=30
environment=DJANGO_SETTINGS_MODULE="fdms_project.settings_production",CELERY_BROKER_URL="redis://127.0.0.1:6379/0",CELERY_RESULT_BACKEND="redis://127.0.0.1:6379/0",PATH="/var/www/fdms/venv/bin:%(ENV_PATH)s"
stdout_logfile=/var/log/fdms/celery_beat.log
stderr_logfile=/var/log/fdms/celery_beat_error.log
stdout_logfile_maxbytes=20MB
stderr_logfile_maxbytes=20MB
redirect_stderr=false
